# infra/edge-agent.service
[Unit]
Description=VisionEdge Edge Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=visionedge
WorkingDirectory=/opt/visionedge/edge-agent
ExecStart=/opt/visionedge/venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
